# k8s/telemetry/fluentd/systemd/out/vmware-whitley/atscale.vmware-test.fluentd.service
[Unit]
Description=Atscale Fluentd Service: vmware-test
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker exec %n stop
ExecStartPre=-/usr/bin/docker rm %n
ExecStartPre=/usr/bin/docker pull prt-registry.sova.intel.com/sraghav1/fluentd-plugins:v3.0.2 
ExecStart=/usr/bin/docker run --name %n \
    --entrypoint=/run.sh \
    -p 0.0.0.0:24211:24211/tcp \
    -v /var/local/common/vmware/systemlogs:/var/local/common/vmware/systemlogs:rw \
    -v /etc/atscale/fluentd/vmware-test:/etc/fluent/config.d:ro \
    --env "FLUENTD_ARGS=-q --no-supervisor" \
    prt-registry.sova.intel.com/sraghav1/fluentd-plugins:v3.0.2 
ExecStop=/usr/bin/docker stop -t 2 %n


[Install]
WantedBy=default.target
