# k8s/telemetry/fluentd/systemd/out/reset-opus/atscale.reset-test.fluentd.service
[Unit]
Description=Atscale Fluentd Service: reset-test
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker exec %n stop
ExecStartPre=-/usr/bin/docker rm %n
ExecStartPre=/usr/bin/docker pull prt-registry.sova.intel.com/sraghav1/fluentd-plugins:v3.0.2 
ExecStart=/usr/bin/docker run --name %n \
    --entrypoint=/run.sh \
    -p 0.0.0.0:24210:24210/tcp \
    -v /var/local/logs/reset_tests:/var/local/logs/reset_tests:rw \
    -v /etc/atscale/fluentd/reset-test:/etc/fluent/config.d:ro \
    --env "FLUENTD_ARGS=-q --no-supervisor" \
    prt-registry.sova.intel.com/sraghav1/fluentd-plugins:v3.0.2 
ExecStop=/usr/bin/docker stop -t 2 %n


[Install]
WantedBy=default.target
